# nginx.conf

description "uWSGI server instance configured to serve app"

start on runlevel [2345]
stop on runlevel [!2345]

setuid user
setgid www-data

env PATH=./env/bin
chdir .
exec uwsgi --ini app.ini

